# ovn_dockerfiles/centos7_ovn_docker/overlay/etc/profile.d/vault.sh
VAULT_ADDR=http://vault:8200
VAULT_TOKEN=vault_token
export VAULT_ADDR
export VAULT_TOKEN
