# workspace/ovn_infrastructure/ansible_automation/roles/jenkins_master/templates/jenkins.service.j2
#######
#!!! Please DO NOT EDIT.  This configuration is maintained by Ansible !!!
#######

[Unit]
Description=Jenkins Daemon

[Service]
Restart=on-failure
Environment="JENKINS_HOME={{ jenkins_home }}/home"

# adding `safeParameters` as we need to whitelist some of the env vars that stash builder is using.
# see https://wiki.jenkins-ci.org/display/JENKINS/Plugins+affected+by+fix+for+SECURITY-170 
ExecStart=/usr/bin/java \
  -XX:+UseG1GC \
  -XX:+ExplicitGCInvokesConcurrent \
  -XX:+ParallelRefProcEnabled \
  -XX:+UseStringDeduplication \
  -XX:+UnlockExperimentalVMOptions \
  -XX:G1NewSizePercent=20 \
  -Dpermissive-script-security.enabled=true \
  -Dhudson.model.UpdateCenter.never=true \
  -Dhudson.model.ParametersAction.safeParameters=pullRequestId \
  -Dcom.sun.jndi.ldap.read.timeout=180000 \
  -Dcom.sun.jndi.ldap.connect.timeout=30000 \
  -jar {{ jenkins_home }}/jenkins.war \
  --handlerCountMax=500 \
  --handlerCountMaxIdle=50 \
  --sessionTimeout=60 \
  --httpListenAddress=127.0.0.1 \
  --httpPort={{ jenkins_port }} \
  --prefix=/{{ jenkins_context_path }} \
  --logfile={{ jenkins_home }}/log/jenkins.log \
  --webroot={{ jenkins_home }}/cache
User={{ jenkins_user }}

[Install]
WantedBy=multi-user.target
