# workspace/ovn_infrastructure/ansible_automation/roles/prometheus/templates/prometheus.service.j2
#####
# PLEASE DO NOT EDIT!!!!  Maintained by Ansible!!!!!
#####
[Unit]
Description=Prometheus is a systems and service monitoring system
Documentation=https://prometheus.io/docs/introduction/overview/

[Service]
Restart=on-failure
User={{ prometheus_user }}
WorkingDirectory={{ prometheus_base_dir }}
ExecStart={{ prometheus_base_dir }}/prometheus \
  -config.file {{ prometheus_base_dir }}/prometheus.yml \
  -storage.local.path {{ prometheus_data_dir }} \
  -storage.local.retention {{ prometheus_data_retention }} \
  -web.route-prefix /{{ prometheus_context_path }} \
  -alertmanager.url http://127.0.0.1:{{ prometheus_alertmanager_port }} \
  -web.listen-address 127.0.0.1:{{ prometheus_port }}
Restart=on-failure

[Install]
WantedBy=multi-user.target
