# workspace/ovn_infrastructure/ansible_roles/prometheus/templates/alertmanager.service.j2
#####
# PLEASE DO NOT EDIT!!!!  Maintained by Ansible!!!!!
#####
[Unit]
Description=The Alertmanager handles alerts sent by client applications such as the Prometheus server
Documentation=https://prometheus.io/docs/alerting/alertmanager/

[Service]
Type=simple
User=prometheus
Group=prometheus
Restart=always
RestartSec=3
TimeoutStopSec=180
StandardOutput=null
StandardError=null
WorkingDirectory={{ alertmanager_base_dir }}
ExecStart={{ alertmanager_base_dir }}/alertmanager \
  -log.format "logger:syslog?appname={{ alertmanager_context_path }}&local=7" \
  -log.level "fatal" \
  -config.file {{ alertmanager_config_dir }}/alertmanager.yml \
  -storage.path {{ alertmanager_install_dir }}/data \
  -data.retention {{ alertmanager_data_retention }} \
  -web.external-url "https://{{ prometheus_external_host }}/{{ alertmanager_context_path }}" \
  -web.listen-address "{{ alertmanager_listen_address }}" \
{%for node in groups.prometheus_alertmanager|default('')%} -mesh.peer "{{ node }}:{{ alertmanager_gossip_port }}"{%endfor%} -mesh.listen-address ":{{ alertmanager_gossip_port }}"
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target