# workspace/ovn_infrastructure/ansible_automation/roles/ovndoc/templates/ovndoc.service.j2
#######
#!!! Please DO NOT EDIT.  This configuration is maintained by Ansible !!!
#######

[Unit]
Description=godoc Daemon

[Service]
Restart=on-failure
Environment="GOPATH={{ ovndoc_home }}"
WorkingDirectory={{ ovndoc_home }}
User={{ ovndoc_user }}

ExecStart=/usr/local/bin/godoc -goroot=src/visa.com/ovn -http=127.0.0.1:{{ ovndoc_port }} -analysis=type -index

[Install]
WantedBy=multi-user.target
